# webui.service
[Unit]
Description=Bus Display Web Configuration Interface
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__HOME__/busdisplay
ExecStart=__HOME__/busdisplay/venv/bin/python __HOME__/busdisplay/webui.py
Restart=always
RestartSec=10
StandardOutput=append:__HOME__/busdisplay/webui.log
StandardError=append:__HOME__/busdisplay/webui.log

[Install]
WantedBy=multi-user.target